# vagrant/hack/discover-worker-nodes.service
[Unit]
After=network.target

[Service]
ExecStart=/home/ubuntu/KubeDeploy/discover-worker-nodes.sh

[Install]
WantedBy=default.target
